# docs/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Example macOS LaunchDaemon for minecrafft-java-console-rx.
  Install:
    sudo cp docs/launchd.plist /Library/LaunchDaemons/com.minecrafft.console.plist
    # edit paths/UserName below, then:
    sudo launchctl bootstrap system /Library/LaunchDaemons/com.minecrafft.console.plist
  Configuration is read from the project's .env file.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.minecrafft.console</string>
  <key>UserName</key><string>localadmin</string>
  <key>WorkingDirectory</key><string>/Users/localadmin/code/minecrafft-java-console-rx</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>src/server.js</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/localadmin/code/minecrafft-java-console-rx/console.log</string>
  <key>StandardErrorPath</key><string>/Users/localadmin/code/minecrafft-java-console-rx/console.log</string>
</dict>
</plist>
